# rockpro64/systemd/ustreamer@.service
# Template service for multi-camera ustreamer instances
# Each instance reads its USB ID and port from /home/noah/plantcam-cam%i.env
# Enable with: sudo systemctl enable --now ustreamer@1 ustreamer@2
# Requires CAM_USB_ID (e.g. 0bda:5844) and CAM_PORT in the env file

[Unit]
Description=ustreamer MJPEG camera server (cam %i)
After=network.target

[Service]
EnvironmentFile=/home/noah/plantcam-cam%i.env
ExecStartPre=/bin/bash -c '/usr/local/bin/find_cameras.sh ${CAM_USB_ID} > /run/ustreamer_device_%i'
ExecStartPre=/bin/bash -c 'v4l2-ctl -d $(cat /run/ustreamer_device_%i) --set-ctrl=auto_exposure=3 --set-ctrl=power_line_frequency=2'
ExecStart=/bin/bash -c 'exec /usr/bin/ustreamer --device $(cat /run/ustreamer_device_%i) --host 0.0.0.0 --port ${CAM_PORT} --resolution 1280x720 --desired-fps 30 --format MJPEG'
TimeoutStartSec=180
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
