# systemd/plantcam-hls.service
# Install on the main server (web host)
# sudo cp plantcam-hls.service /etc/systemd/system/
# sudo systemctl daemon-reload && sudo systemctl enable --now plantcam-hls
# Requires: /home/noah/plantcam.env with ROCKPRO64_IP, CAM_PORT, HLS_DIR

[Unit]
Description=Plant Cam HLS transcoder (MJPEG to H.264 VAAPI)
After=network.target

[Service]
EnvironmentFile=/home/noah/plantcam.env
ExecStartPre=/bin/mkdir -p ${HLS_DIR}
ExecStart=/usr/bin/ffmpeg \
    -use_wallclock_as_timestamps 1 \
    -fflags +genpts+discardcorrupt \
    -reconnect 1 -reconnect_streamed 1 -reconnect_delay_max 2 \
    -i http://${ROCKPRO64_IP}:${CAM_PORT}/?action=stream \
    -r 15 \
    -vaapi_device /dev/dri/renderD128 \
    -vf "lagfun=decay=0.9995,tmedian=radius=4:percentile=1,eq=brightness=-0.08:saturation=0.85:contrast=1.1,colorbalance=rm=-0.15:gm=-0.05:bm=0.15:rh=-0.1:bh=0.1,drawtext=text='%%{localtime}':fontsize=20:fontcolor=white:borderw=2:bordercolor=black:x=10:y=h-35,format=nv12,hwupload" \
    -c:v h264_vaapi \
    -b:v 1500k \
    -g 30 \
    -f hls \
    -hls_time 2 \
    -hls_list_size 30 \
    -hls_flags delete_segments+temp_file \
    -hls_segment_filename ${HLS_DIR}/segment_%%03d.ts \
    ${HLS_DIR}/stream.m3u8
Restart=always
RestartSec=10
User=noah

[Install]
WantedBy=multi-user.target
